# scripts/init/openrc.sh
#!/sbin/openrc-run

description="placeholder_kernel_desc"
command="placeholder_cmd_path"
command_args="placeholder_cmd_args"
pidfile="placeholder_pid_path"
output_log="placeholder_log_path"
error_log="placeholder_log_path"
command_background=true
