# script/init/OpenRC.sh
#!/sbin/openrc-run
# shellcheck disable=SC2034

# 服务元数据
description="placeholder_kernel_desc"
command="placeholder_cmd_path"
command_args="placeholder_cmd_args"
pidfile="placeholder_pid_file"
output_log="placeholder_log_file"
error_log="placeholder_log_file"
command_background=true

# 依赖关系
depend() {
    after network localmount
    need net
}


